# 16/provision/client/borg-backup.timer
[Unit]
Description=Borg Backup

[Timer]
OnUnitActiveSec=5min
Unit=borgservice.service

[Install]
WantedBy=timers.target
